# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; ── 1. Xvfb 虚拟显示器 ────────────────────────────────────
; 创建一个虚拟屏幕，Chrome 在上面渲染画面
[program:xvfb]
command=Xvfb :%(ENV_DISPLAY_NUM)s -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── 2. fluxbox 窗口管理器 ─────────────────────────────────
; 让 Chrome 窗口能正常显示（不用窗口管理器 Chrome 可能无法渲染）
[program:fluxbox]
command=fluxbox
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s"

; ── 3. x11vnc VNC 服务器 ──────────────────────────────────
; 把虚拟显示器的画面通过 VNC 协议暴露出去
[program:x11vnc]
command=/bin/sh -c "if [ -n \"$VNC_PASSWORD\" ]; then x11vnc -storepasswd \"$VNC_PASSWORD\" /tmp/vncpass && exec x11vnc -display :%(ENV_DISPLAY_NUM)s -forever -shared -rfbauth /tmp/vncpass -rfbport %(ENV_VNC_PORT)s -noxdamage -wait 5; else exec x11vnc -display :%(ENV_DISPLAY_NUM)s -forever -shared -nopw -rfbport %(ENV_VNC_PORT)s -noxdamage -wait 5; fi"
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

; ── 4. noVNC 网页 VNC 客户端 ──────────────────────────────
; 把 VNC 转成 WebSocket，让浏览器可以直接访问
[program:novnc]
command=websockify --web /usr/share/novnc %(ENV_NOVNC_PORT)s 127.0.0.1:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40

; ── 5. Chromium（有界面模式）─────────────────────────────
; 运行在虚拟显示器上，开启 CDP 远程调试
[program:chromium]
command=chromium --no-sandbox --disable-dev-shm-usage --remote-debugging-address=127.0.0.1 --remote-debugging-port=%(ENV_BB_CDP_PORT)s --user-data-dir=/data/chrome-profile --no-first-run --no-default-browser-check --disable-translate --window-size=%(ENV_SCREEN_WIDTH)s,%(ENV_SCREEN_HEIGHT)s --start-maximized about:blank
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s"

; ── 6. bb-browser daemon ─────────────────────────────────
[program:bb-daemon]
command=node /app/dist/daemon.js --host 0.0.0.0 --port %(ENV_BB_DAEMON_PORT)s --cdp-host 127.0.0.1 --cdp-port %(ENV_BB_CDP_PORT)s
autostart=true
autorestart=true
startretries=5
startsecs=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=60
environment=BB_BROWSER_HOME="/data/bb-browser"
