# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; ── Chromium（无头模式，提供 CDP 接口）────────────────────
[program:chromium]
command=chromium
    --no-sandbox
    --disable-dev-shm-usage
    --disable-gpu
    --headless=new
    --remote-debugging-address=127.0.0.1
    --remote-debugging-port=%(ENV_BB_CDP_PORT)s
    --user-data-dir=/data/chrome-profile
    --disable-extensions
    --disable-background-networking
    --disable-default-apps
    --disable-sync
    --no-first-run
    --no-default-browser-check
    --disable-translate
    --disable-features=TranslateUI
    about:blank
autostart=true
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── bb-browser daemon ─────────────────────────────────────
[program:bb-daemon]
command=node /app/dist/daemon.js
    --host 0.0.0.0
    --port %(ENV_BB_DAEMON_PORT)s
    --cdp-host 127.0.0.1
    --cdp-port %(ENV_BB_CDP_PORT)s
autostart=true
autorestart=true
startretries=5
startsecs=5
; 等待 Chromium 先启动
depends_on=chromium
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
environment=BB_BROWSER_HOME="/data/bb-browser"
